# stra2us/backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server --dir /data --appendonly yes
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.out.log
stderr_logfile=/var/log/supervisor/redis.err.log

[program:stra2us]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 8153
directory=/app/src
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/stra2us.out.log
stderr_logfile=/var/log/supervisor/stra2us.err.log
environment=PYTHONPATH="/app/src"
